# crontab
@reboot rm -rf /home/pi/livestream/timelapse
@reboot /usr/local/bin/node /home/pi/livestream/index.js
