# crontab
@reboot rm -rf /home/pi/livestream/timelapse
@reboot sleep 60 && /usr/local/bin/node /home/pi/livestream/index.js
